# cashex-bot.service
[Unit]
Description=CASHEX Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cashex
ExecStart=/opt/cashex/venv/bin/python /opt/cashex/bot.py
Restart=always
RestartSec=10
EnvironmentFile=/opt/cashex/.env

[Install]
WantedBy=multi-user.target
